# deploy/keepalive.timer
[Unit]
Description=Ping aliexpress-oauth.onrender.com every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true
Unit=aliexpress-oauth-keepalive.service

[Install]
WantedBy=timers.target
